# provisioning/docker-files/php_fpm/supervisord.conf
[supervisord]
nodaemon=true

[program:php-fpm]
logfile=/var/log/supervisord/php_fpm.log
pidfile=/var/log/supervisord/php_fpm.pid

command=php-fpm -F
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autorestart=true
startretries=0